# arch/systemd/ulauncher.service.d/gnome-only.conf
[Unit]
ConditionEnvironment=XDG_CURRENT_DESKTOP=GNOME
